# ubicom-linux-dist-3.2.1/openwrt/target/linux/ubicom32/IP8100RGW/base-files/etc/init.d/ifplugd
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 Ubicom, Inc

START=99

start() {
	local wan_device
	local wanmon_enabled

	config_load network
	config_get wanmon_enabled wanmon enabled
	if [ "$wanmon_enabled" == "0" ]; then
		echo "wanmon not enabled. Do not run ifplugd."
		exit 0
	fi

	config_get wan_device wan device

	if [ -n "$wan_device" -a -n "$(which ifplugd)" ]; then
		ifplugd -i $wan_device -p -q -u 0 -d 0 -I -m iff
	fi
}

stop() {
	killall ifplugd
}
